# dev-db/ctdb/files/ctdb.initd-r1
#!/sbin/runscript

command="/usr/sbin/ctdbd_wrapper"

depend() {
  need localmount
  need net
  before samba
  after bootmisc
}

start_pre() {
	source /etc/ctdb/ctdbd.conf || return 1

	if [ -z "${CTDB_PUBLIC_ADDRESSES}" ] ; then
		eerror "You must configure the location of the CTDB_PUBLIC_ADDRESSES in /etc/ctdb/ctdbd.conf"
		return 1
	fi

	# check all persistent databases that they look ok
	PERSISTENT_DB_DIR="/var/lib/ctdb/persistent"
	if [ -n "${CTDB_DBDIR}" ] ; then
		PERSISTENT_DB_DIR="${CTDB_DBDIR}/persistent"
	fi

	checkpath -q -d ${PERSISTENT_DB_DIR}

	for PDBASE in $(ls ${PERSISTENT_DB_DIR}/*.tdb.[0-9] 2>/dev/null) ; do
		if ! /usr/bin/tdbdump ${PDBASE} >/dev/null 2>/dev/null ; then
			eerror "Persistent database ${PDBASE} is corrupted! CTDB will not start."
			return 1
		fi
	done
}

start() {
	start-stop-daemon --start --background --exec ${command} -- start
}

stop() {
	start-stop-daemon --start --exec ${command} --pidfile ${CTDB_PIDFILE} \
		-- stop \
		&& mark_service_stopped
}
