# net-analyzer/arpalert/files/arpalert.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# Written by Polynomial-C

pidfile="/var/run/arpalert.lock"
command="/usr/sbin/arpalert"
command_args="-d"

depend() {
	need net
}
