# app-laptop/dell-fan-mon/files/dell-fan-mon.init
#!/sbin/openrc-run

description="Monitors fans of Dell Laptops"
pidfile="/run/dell-fan-mon.pid"
command="/usr/sbin/dell-fan-mon"

depend() {
	after modules
}

start_pre() {
	if [ ! -e /proc/i8k ] ; then
		eerror "The i8k driver is not installed"
		return 1
	fi
}
